# conf/supervise-honeypress.conf
[program:honeypress]
command=/opt/honeypress/honeypress.py
redirect_stderr=true
stdout_logfile=/opt/honeypress/logs/access.log
autorestart=true
